# templates/dist/etc/cron.d/flight-direct.erb
# Flight Direct cron tab file for setting up reboot scripts
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=<%= ENV['USER'] %>
@reboot <%= ENV['USER'] %> run-parts <%= ENV['FL_ROOT'] %>/etc/cron.reboot
