# etc/profile.d/01-banner.csh
#!/bin/csh
/bin/bash $FL_ROOT/etc/profile.d/01-banner.sh
