# templates/dist/etc/profile.d/flight-direct.csh.erb
################################################################################
##
## Flight Direct - profile configuration file
## Copyright (c) 2018 Alces Software Ltd
##
################################################################################
setenv FL_ROOT <%= ENV['FL_ROOT'] %>
source $FL_ROOT/etc/profile.csh
