# templates/dist/etc/profile.d/flight-direct.sh.erb
################################################################################
##
## Flight Direct - profile configuration file
## Copyright (c) 2018 Alces Software Ltd
##
################################################################################
source <%= ENV['FL_ROOT'] %>/etc/profile.sh
